# supervisor/n2kafka.conf
[program:n2kafka]
command='./n2kafka config.json'
autostart=true
autorestart=true
