# supervisor/n2kafka.ini
[program:n2kafka]
command=/app/n2kafka /app/config.json
autostart=true
autorestart=true
stdout_logfile=/var/log/n2kafka.log
stderr_logfile=/var/log/n2kafka-err.log
